# systemd/sammurai-digest.service
[Unit]
Description=Sammurai Daily Digest and Wiki Compilation
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/chhikv/ai/sammurai
ExecStart=/home/chhikv/.venv/bin/python -m src.backend.digest_runner
EnvironmentFile=%h/.config/sammurai/env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
